# daily_trello.timer
[Unit]
Description=Runs trello collector once a day

[Timer]
OnCalendar=daily
Unit=daily_trello.service

[Install]
WantedBy=multi-user.target